# deployment/gunicorn/gunicorn_config.ini
[unix_http_socket]
file = /var/run/gunicorn/gunicorn.sock

[inet]
host = 127.0.0.1
port = 8000

[program:gunicorn]
command = gunicorn config.wsgi:application --workers 4 --worker-class sync --bind unix:/var/run/gunicorn/gunicorn.sock --error-logfile /var/log/gunicorn/error.log --access-logfile /var/log/gunicorn/access.log
directory = /app
user = www-data
autostart = true
autorestart = true
